# PycharmProjects_2017/python3/flask_1.conf
start on [2345]
stop on [06]

script
    exec uwsgi --ini ./flask_1.ini
end script

